# systemd/charge-limiter.service
[Unit]
Description=Battery charge limiter at 50%

[Service]
Type=oneshot
ExecStart=/usr/local/bin/charge-limiter.sh
